# init/centos/7.1x64/etc/profile.d/java.sh
#!/bin/bash
JAVA_HOME=/usr/java/default
M2_HOME=/opt/maven/apache-maven
PATH=$JAVA_HOME/bin:$PATH
export PATH M2_HOME JAVA_HOME
#export CLASSPATH=.
