# CentOS/etc/profile.d/java.sh
#!/bin/bash
JAVA_HOME=/usr/java/jdk1.8.0_25
M2_HOME=/opt/maven/apache-maven
PATH=$JAVA_HOME/bin:$PATH
export PATH M2_HOME JAVA_HOME
#export CLASSPATH=.
